# venus/services/bluetti-repair-on-boot/run
#!/bin/sh
set -eu

APP_DIR="${BLUETTI_GATEWAY_APP_DIR:-/data/bluetti-venus-gateway}"
LOG_DIR="${BLUETTI_DATA_DIR:-/data/bluetti-gateway}/logs"
LOG_FILE="$LOG_DIR/bluetti-repair-on-boot.log"

mkdir -p "$LOG_DIR"
"$APP_DIR/venus/repair-if-needed.sh" >>"$LOG_FILE" 2>&1 || true

# runit restarts services that exit. Keep this boot-time repair service
# alive after the one-shot repair so Venus does not run repair in a loop.
exec sleep 2147483647
