# venus/services/bluetti-collector/run
#!/bin/sh
set -eu

APP_DIR="${BLUETTI_GATEWAY_APP_DIR:-/data/bluetti-venus-gateway}"
DATA_DIR="${BLUETTI_DATA_DIR:-/data/bluetti-gateway}"
CONFIG_FILE="${BLUETTI_CONFIG_FILE:-$DATA_DIR/bluetti-gateway.env}"
LOG_DIR="$DATA_DIR/logs"
LOG_FILE="$LOG_DIR/bluetti-collector.log"

mkdir -p "$LOG_DIR" "$DATA_DIR/cache" "$DATA_DIR/certs" /run/bluetti-gateway

export PYTHONPATH="$APP_DIR/src:${PYTHONPATH:-}"
export BLUETTI_DATA_DIR="$DATA_DIR"
export BLUETTI_RUN_DIR="${BLUETTI_RUN_DIR:-/run/bluetti-gateway}"

exec python3 -m bluetti_venus_gateway.services.collector --config "$CONFIG_FILE" >>"$LOG_FILE" 2>&1
